# postgresql.upstart
description "Starts PostgreSQL instances for all configured clusters"

start on started SERVICE=syslog

env PG_DB_CLUSTERS

pre-start script
	exec >/tmp/pre-start 2>&1
	set -x
	if [ -z "$PG_DB_CLUSTERS" ] ; then
		# read PG_DB_CLUSTERS from /etc/sysconfig/postgresql
		# for init script compatibility
		. /etc/sysconfig/postgresql
	fi
	for DATADIR in $PG_DB_CLUSTERS ; do
		start postgresql/instance DATADIR=$DATADIR
	done
end script
