# postgresql.target
[Unit]
Description=PostgreSQL server instances
